# deploy/systemd/bandwidth-analyzer.timer
[Unit]
Description=Run bandwidth analyzer daily at 23:59

[Timer]
OnCalendar=*-*-* 23:59:00
Persistent=true

[Install]
WantedBy=timers.target
