# monitor/systemd/bandwidth-daily-report.timer
[Unit]
Description=Send daily traffic report at 23:00

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
