# monitor/systemd/route-monitor.service
[Unit]
Description=Dynamic Route Change Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/bandwidth_monitor/route_daemon.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=route-monitor

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/log/bandwidth
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
